# debian/crontab
# Check if tcpquotad is up and running...
15 * * * *      /usr/lib/tcpquota/verifytcpquotad
